# sd64/etc/xinetd.d/qmserver
# This service supports connections to the SD
# database server from telnet clients.
# If your SD binary is located in a place other
# than /usr/local/sdsys/bin, you'll need to update the server
# entry below to reflect that new path.
#
# Port #4242 is normally used for this.

service qmserver
{
	protocol 	= tcp
	flags		= REUSE
	socket_type	= stream
	wait		= no
	groups		= yes
	user		= root
	group		= sdusers
	umask		= 002
	server		= /usr/local/sdsys/bin/sd
	server_args	= -n
	log_on_failure	+= USERID
	disable		= no
}
