# sd64/etc/xinetd.d/services
# entries for your /etc/services file.
sdserver	4244/tcp	# SDServer
sdclient	4245/tcp	# SDClient
