# sd64/usr/lib/systemd/system/sdclient@.service
[Unit]
Description=String Database (SD)  Per-Connection clinet Server

[Service]
ExecStart=/usr/local/sdsys/bin/sd -n -q
StandardInput=socket

